# deployer/deployer.service
[Unit]
Description=Deploy webhook server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/mnt/services/deployer/deployer
Restart=on-failure
RestartSec=5s

# Environment
Environment=PORT=3201
Environment=DEPLOY_PATH=/mnt/services
EnvironmentFile=-/etc/deployer/env

# Security hardening
User=rnixon
Group=docker
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/mnt/services
PrivateTmp=true

[Install]
WantedBy=multi-user.target
