# no-time-to-explain/systemd/no-time-to-explain-refresh.service
[Unit]
Description=No Time to Explain - Refresh Data
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/mnt/services/no-time-to-explain
ExecStart=/usr/bin/docker compose run --rm app /app/no-time-to-explain refresh

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=no-time-to-explain-refresh

[Install]
WantedBy=multi-user.target
